# config/kanata/.config/systemd/user/kanata.service
[Unit]
Description=Kanata keyboard remapper
Documentation=https://github.com/jtroo/kanata
After=graphical-session.target

[Service]
Type=simple
ExecStart=%h/.local/bin/kanata -c %h/.config/kanata/kanata.kbd -n
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
